# debian/tatu-repo-devel.cron.hourly
#!/bin/sh

. /etc/taturepo.conf

if [ -f "/etc/awstats/awstats.${REPO_SERVER}.conf" ]; then
    /usr/lib/cgi-bin/awstats.pl -config=${REPO_SERVER} -update > /dev/null
fi
